# Hallo/smart-srun-main/root/etc/init.d/smart_srun
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

SERVICE_NAME="smart_srun"
PROG="/usr/lib/smart_srun/client.py"

find_python() {
	command -v python3 2>/dev/null || command -v python3.11 2>/dev/null
}

start_service() {
	PYTHON_BIN="$(find_python)"
	if [ -z "$PYTHON_BIN" ]; then
		logger -t "$SERVICE_NAME" "python3 not found, please install python3-light"
		return 1
	fi

	procd_open_instance
	procd_set_param command "$PYTHON_BIN" -B "$PROG" daemon
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

stop_client_processes() {
	for proc in /proc/[0-9]*; do
		[ -r "$proc/cmdline" ] || continue
		cmdline="$(tr '\000' ' ' < "$proc/cmdline" 2>/dev/null)"
		case "$cmdline" in
			*"$PROG"*)
				kill -TERM "${proc##*/}" 2>/dev/null
				;;
		esac
	done
	sleep 1
	for proc in /proc/[0-9]*; do
		[ -r "$proc/cmdline" ] || continue
		cmdline="$(tr '\000' ' ' < "$proc/cmdline" 2>/dev/null)"
		case "$cmdline" in
			*"$PROG"*)
				kill -KILL "${proc##*/}" 2>/dev/null
				;;
		esac
	done
}

stop_service() {
	stop_client_processes
	return 0
}
