# Hallo/luci-app-autoweblogin/etc/init.d/autoweblogin
#!/bin/sh /etc/rc.common  

START=99  

log_file="/tmp/log/autoweblogin.log"  
ENABLED="$(uci get autoweblogin.config.enabled)"  

start() {  
  if [ "$ENABLED" -eq 0 ]; then  
    echo "[$(date "+%Y-%m-%d %H:%M:%S")]: 已停用" >> "$log_file"  
    return 1  
  fi  

  # 检查进程是否运行  
  if pgrep -f "sh /etc/init.d/autoweblogin" > /dev/null; then  
    echo "[$(date "+%Y-%m-%d %H:%M:%S")]: 认证进程已运行" >> "$log_file"  
    return  
  fi  

  sh /etc/init.d/autoweblogin &  
  echo "[$(date "+%Y-%m-%d %H:%M:%S")]: 启用认证" >> "$log_file"  
}  

stop() {  
  pid=$(pgrep -f "sh /etc/init.d/autoweblogin")  
  if [ -n "$pid" ]; then  
    kill "$pid" 2>/dev/null  
    echo "[$(date "+%Y-%m-%d %H:%M:%S")]: 已停止认证" >> "$log_file"  
  else  
    echo "[$(date "+%Y-%m-%d %H:%M:%S")]: 认证进程未运行" >> "$log_file"  
  fi  
}  

restart() {  
  stop  
  sleep 1  
  start  
  echo "[$(date "+%Y-%m-%d %H:%M:%S")]: 已重新启动" >> "$log_file"  
}
